# openhab/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: openHab
# Runs the Nginx daemon
# ==============================================================================

# Use server mode with the default command when there is no pseudo-TTY
#if [ "$interactive" == "false" ] && [ "$(IFS=" "; echo "$@")" == "gosu openhab tini -s ./start.sh" ]; then
#    command=($@ server)
#    exec "${command[@]}"
#else
#    exec "$@"
#fi

# Wait for grafana to become available
bashio::net.wait_for 8080

bashio::log.info "Starting NGinx..."

exec nginx